# packaging/pacman/lb2d.install
post_install() {
    systemctl daemon-reload
    systemctl enable lb2d.service
    systemctl start lb2d.service
}

post_upgrade() {
    systemctl daemon-reload
    if [[ "$(systemctl is-enabled lb2d.service)" = "enabled" ]]; then
        systemctl restart lb2d.service
    fi
}